# udev.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVERS=="usb", RUN+="/usr/bin/carusbam $attr{busnum} $attr{devnum}"